# aurral-search/systemd/aurral-search-update.service
[Unit]
Description=Refresh Aurral Search canonical index
After=network.target docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/aurral-search
ExecStart=/opt/aurral-search/update.sh
